# deploy/systemd/regknots-worker.service.d/memory.conf
# systemd drop-in: cgroup memory caps for regknots-worker (Celery).
#
# Steady-state ~110 MB. Caps mirror regknots-web. Celery beat tasks that
# legitimately need more memory should be moved to scripts/run_ingest.sh
# (their own cgroup, 1.5 GB cap) — not allowed to bloat the worker.

[Service]
MemoryHigh=512M
MemoryMax=1G
