# infra/regknots-worker.service
[Unit]
Description=RegKnots Celery Worker + Beat
After=network.target redis.service postgresql.service
Requires=redis.service

[Service]
Type=simple
User=root
WorkingDirectory=/opt/RegKnots/apps/api
Environment=PATH=/root/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
EnvironmentFile=/opt/RegKnots/.env
ExecStart=/root/.local/bin/uv run celery -A celery_beat worker --beat --loglevel=info --concurrency=2
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=regknots-worker

[Install]
WantedBy=multi-user.target
