# deploy/systemd/regknots-backup.service
[Unit]
Description=RegKnots daily Postgres backup — pg_dump from regknots-postgres container, gzipped, 14-day retention
Documentation=https://github.com/blakemarchal/RegKnots
# Don't try to back up if the database isn't running.
After=docker.service
Requires=docker.service

[Service]
Type=oneshot
WorkingDirectory=/opt/RegKnots
ExecStart=/opt/RegKnots/scripts/backup_postgres.sh
StandardOutput=journal
StandardError=journal
# Headroom for the pg_dump → pigz -3 pipe. Empirically the dump is ~1.5GB
# uncompressed → ~400MB compressed; pigz -3 on 2 cores runs in ~3-5min on
# this VPS. 30min cap is comfortable; if a backup ever runs past it,
# something else is wrong (CPU contention, disk full, DB locked).
TimeoutStartSec=1800

[Install]
WantedBy=multi-user.target
