# deploy/systemd/regknots-web.service.d/memory.conf
# systemd drop-in: cgroup memory caps for regknots-web (Next.js).
#
# Tighter caps than the api because the Next.js prod server has a smaller
# steady-state memory footprint (~250 MB observed) and shouldn't ever need
# to grow past 1 GB. If it does, something is wrong (memory leak in a
# server component, runaway prefetch, etc.) and a restart is preferable
# to taking the whole box down.

[Service]
MemoryHigh=512M
MemoryMax=1G
