# deploy/systemd/regknots-api.service.d/memory.conf
# systemd drop-in: cgroup memory caps for regknots-api.
#
# Sprint post-D6.83 audit (2026-05-08). The audit found 116 service
# failures over 14 days, including one explicit OOM-kill on May 1 where
# a uvicorn worker grew to 2.4 GB anon-rss and triggered global OOM —
# taking the box (including SpiritFlow co-tenant) down for 9 hours.
#
# With these caps:
#   MemoryHigh=1.5G — kernel begins throttling allocations + reclaiming
#                     pages above this. Soft pressure, no kill.
#   MemoryMax=2G    — hard kill: cgroup OOM-kills the api process at 2G,
#                     systemd auto-restarts via Restart=always, the rest
#                     of the box stays healthy.
#
# Drop-in pattern means this overlays the master regknots-api.service
# unit without modifying it. Revert by deleting this file and running
# `systemctl daemon-reload && systemctl restart regknots-api`.

[Service]
MemoryHigh=1.5G
MemoryMax=2G
